# roles/setup/templates/etc/profile.d/gradle.sh.j2
export PATH=$PATH:/opt/gradle-5.3.1/bin
export GRADLE_HOME=/opt/gradle-5.3.1
