# roles/setup/templates/etc/profile.d/groovy.sh.j2
export GROOVY_HOME="/opt/groovy-{{ groovy_version }}/"
export PATH="$PATH:$GROOVY_HOME/bin"
